# configs/label-printer.service
[Unit]
Description=QuickBooks Label Printer
After=network.target

[Service]
Type=simple
User=labelprinter
Group=labelprinter
WorkingDirectory=/opt/label-printer
Environment=PATH=/opt/label-printer/venv/bin
Environment=FLASK_ENV=production
Environment=FLASK_APP=app.py
ExecStart=/opt/label-printer/venv/bin/gunicorn --config configs/gunicorn.conf.py app:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/label-printer

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target